# smeserver-nfs-1.0.0-1-filelist
%dir %attr(0755,root,root) /etc/e-smith
%dir %attr(0755,root,root) /etc/e-smith/db
%dir %attr(0755,root,root) /etc/e-smith/db/configuration
%dir %attr(0755,root,root) /etc/e-smith/db/configuration/defaults
%dir %attr(0755,root,root) /etc/e-smith/db/configuration/defaults/nfslock
/etc/e-smith/db/configuration/defaults/nfslock/type
/etc/e-smith/db/configuration/defaults/nfslock/status
/etc/e-smith/db/configuration/defaults/nfslock/TCPPorts
/etc/e-smith/db/configuration/defaults/nfslock/access
/etc/e-smith/db/configuration/defaults/nfslock/UDPPorts
/etc/e-smith/db/configuration/defaults/nfslock/statdOutgoingPort
/etc/e-smith/db/configuration/defaults/nfslock/lockdPort
/etc/e-smith/db/configuration/defaults/nfslock/statdPort
%dir %attr(0755,root,root) /etc/e-smith/db/configuration/defaults/nfs
/etc/e-smith/db/configuration/defaults/nfs/type
/etc/e-smith/db/configuration/defaults/nfs/status
/etc/e-smith/db/configuration/defaults/nfs/TCPPorts
/etc/e-smith/db/configuration/defaults/nfs/UDPPorts
/etc/e-smith/db/configuration/defaults/nfs/access
/etc/e-smith/db/configuration/defaults/nfs/mountdPort
/etc/e-smith/db/configuration/defaults/nfs/rquotadPort
%dir %attr(0755,root,root) /etc/e-smith/db/configuration/defaults/portmap
/etc/e-smith/db/configuration/defaults/portmap/type
/etc/e-smith/db/configuration/defaults/portmap/status
/etc/e-smith/db/configuration/defaults/portmap/UDPPort
/etc/e-smith/db/configuration/defaults/portmap/TCPPort
/etc/e-smith/db/configuration/defaults/portmap/access
%dir %attr(0755,root,root) /etc/e-smith/templates
%dir %attr(0755,root,root) /etc/e-smith/templates/etc
%dir %attr(0755,root,root) /etc/e-smith/templates/etc/hosts.allow
/etc/e-smith/templates/etc/hosts.allow/lockd
/etc/e-smith/templates/etc/hosts.allow/mountd
/etc/e-smith/templates/etc/hosts.allow/portmap
/etc/e-smith/templates/etc/hosts.allow/rquotad
/etc/e-smith/templates/etc/hosts.allow/statd
%dir %attr(0755,root,root) /etc/e-smith/templates/etc/sysconfig
%dir %attr(0755,root,root) /etc/e-smith/templates/etc/sysconfig/nfs
/etc/e-smith/templates/etc/sysconfig/nfs/statdPort
/etc/e-smith/templates/etc/sysconfig/nfs/lockdPort
/etc/e-smith/templates/etc/sysconfig/nfs/mountdPort
/etc/e-smith/templates/etc/sysconfig/nfs/rquotadPort
%dir %attr(0755,root,root) /etc/e-smith/events
%dir %attr(0755,root,root) /etc/e-smith/events/console-save
%dir %attr(0755,root,root) /etc/e-smith/events/console-save/templates2expand
%dir %attr(0755,root,root) /etc/e-smith/events/console-save/templates2expand/etc
/etc/e-smith/events/console-save/templates2expand/etc/hosts.allow
%dir %attr(0755,root,root) /etc/e-smith/events/console-save/templates2expand/etc/sysconfig
/etc/e-smith/events/console-save/templates2expand/etc/sysconfig/nfs
%dir %attr(0755,root,root) /etc/e-smith/events/bootstrap-console-save
%dir %attr(0755,root,root) /etc/e-smith/events/bootstrap-console-save/templates2expand
%dir %attr(0755,root,root) /etc/e-smith/events/bootstrap-console-save/templates2expand/etc
/etc/e-smith/events/bootstrap-console-save/templates2expand/etc/hosts.allow
%dir %attr(0755,root,root) /etc/e-smith/events/bootstrap-console-save/templates2expand/etc/sysconfig
/etc/e-smith/events/bootstrap-console-save/templates2expand/etc/sysconfig/nfs
%dir %attr(0755,root,root) /etc/e-smith/events/remoteaccess-update
%dir %attr(0755,root,root) /etc/e-smith/events/remoteaccess-update/templates2expand
%dir %attr(0755,root,root) /etc/e-smith/events/remoteaccess-update/templates2expand/etc
/etc/e-smith/events/remoteaccess-update/templates2expand/etc/hosts.allow
%dir %attr(0755,root,root) /etc/e-smith/events/remoteaccess-update/templates2expand/etc/sysconfig
/etc/e-smith/events/remoteaccess-update/templates2expand/etc/sysconfig/nfs
%dir %attr(0755,root,root) /etc/rc.d/rc7.d
/etc/rc.d/rc7.d/S13portmap
/etc/rc.d/rc7.d/S59nfslock
/etc/rc.d/rc7.d/S60nfs
%doc COPYING
